# deploy/systemd/lss-headscale-dashboard.service
[Unit]
Description=LSS Headscale Dashboard
After=network.target headscale.service
Wants=network.target

[Service]
Type=simple
User=lss-dashboard
Group=lss-dashboard
ExecStart=/usr/local/bin/lss-headscale-dashboard --config /etc/lss-headscale-dashboard/config.yaml
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ReadWritePaths=/var/lib/lss-headscale-dashboard
# Uncomment ONLY if installing on the same host as Headscale (so the service
# can read /var/run/headscale/headscale.sock). systemd will refuse to start
# the unit (status=216/GROUP) if this group does not exist.
# SupplementaryGroups=headscale

[Install]
WantedBy=multi-user.target
